# scripts/systemd/telegram-shop-backup.timer
[Unit]
Description=Run Telegram Shop Platform production backup daily at 04:00 Moscow time

[Timer]
OnCalendar=*-*-* 04:00:00 Europe/Moscow
AccuracySec=5min
Persistent=true
Unit=telegram-shop-backup.service

[Install]
WantedBy=timers.target
